# crontab
* * * * * . /etc/profile;/bin/sh /root/server/monitoring.sh

crontab -e
crontab -l
